# service/barrier_stop.service
[Unit]
Description=Stop Barrier
Requires=barrier.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop barrier.service

