# service/barrier_stop.timer
[Unit]
Description=Timer to stop Barrier service

[Timer]
OnCalendar=Mon-Fri 17:00:00
Unit=barrier_stop.service

[Install]
WantedBy=timers.target
