# service/barrier.timer
[Unit]
Description=Timer to start Barrier service

[Timer]
OnCalendar=weekday 08:00:00

[Install]
WantedBy=timers.target
